# src/lic/utils/sdev3.crontab
2,7,12,17,22,27,32,37,42,47,52,57 * * * * /bt/ab/util/sdev3-often 1> /bt/ab/logs/sdev3-often.log 2>&1
